# recipes-extended/usbip-server/usbip-server/usbipd.service
[Unit]
Description=USBIP server service

[Service]
ExecStart=usbipd
