# recipes-extended/usbip-server/usbip-server/usbip-share@.service
[Unit]
Description=USBIP service to share an USB device through network
Wants=usbipd.service

[Service]
Type=oneshot
ExecStart=usbip bind -b %i
